# devops/roles/app_setup/templates/supervisor.ini
[program:{{ app_name }}]
command=/usr/local/bin/uwsgi
  --socket /tmp/{{ app_name }}.sock
  --logto {{ webapps_dir }}/{{ app_name }}/log/uwsgi/{{ app_name }}_uwsgi.log
  --home /home/{{ deploy_user }}/.virtualenvs/{{ app_name }}
  --pythonpath {{ webapps_dir }}/{{ app_name }}/html
  --wsgi-file {{ webapps_dir }}/{{ app_name }}/{{ app_root }}/wsgi.py
  --callable {{ wsgi_callable }}
  --max-requests 1000
  --master
  --processes 1
  --chmod
  --uid {{ deploy_user }}
directory={{ webapps_dir }}/{{ app_name }}/{{ app_root }}
logfile={{ webapps_dir }}/{{ app_name }}/log/supervisor/{{ app_name }}_supervisor.log
autostart=true
autorestart=true
