# devops/roles/supervisord_setup/templates/supervisor.ini
[group:{{ app_name }}]
programs=touchlog,logpermissions,application
priority=999

[program:touchlog]
command=/usr/bin/touch {{ webapps_dir }}/{{ app_name }}/logs/ecodriving_application.log
autorestart=false

[program:logpermissions]
command=/bin/chown {{ deploy_username }}:{{ deploy_username }} {{ webapps_dir }}/{{ app_name }}/logs/ecodriving_application.log
autorestart=false

[program:application]
command=/usr/local/bin/uwsgi
  --socket /tmp/{{ app_name }}.sock
  --logto {{ webapps_dir }}/{{ app_name }}/logs/uwsgi/kyburz_uwsgi_application.log
  --home {{ webapps_dir }}/{{ app_name }}/venv
  --pythonpath {{ webapps_dir }}/{{ app_name }}/src/Kyburz_flask
  --wsgi-file {{ webapps_dir }}/{{ app_name }}/src/Kyburz_flask/{{ wsgi_file }}
  --callable {{ wsgi_callable }}
  --max-requests 1000
  --master
  --processes 1
  --chmod
  --uid deploy
directory={{ webapps_dir }}/{{ app_name }}/src
logfile={{ webapps_dir }}/{{ app_name }}/logs/supervisor/kyburz_supervisor.log
autostart=true
autorestart=true
